# sshd.conf
[program:sshd]
command = /usr/sbin/sshd -D
process_name = sshd
