# test/scripts/service/xdummy/run
#!/bin/sh

exec 2>&1
exec /usr/bin/Xorg \
		-noreset -nolisten tcp \
		+extension GLX +extension RANDR +extension RENDER \
		-logfile /var/log/xdummy.log \
		:0
