# test/scripts/service/dbus-daemon/run
#!/bin/sh

exec 2>&1
exec /usr/bin/chpst -u messagebus /usr/bin/dbus-daemon --system --nofork
